# laravel/supervisord.conf
# /Users/medaminebentaieb/Docker_files/Laravel_Docker/laravel/supervisord.conf

[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.err.log
